# deploy/systemd/vefr.service
[Unit]
Description=vefr IPv6 HTTP forward proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vefr
Group=vefr
ExecStart=/usr/local/bin/vefr -config /etc/vefr/config.toml
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
MemoryDenyWriteExecute=true
ReadOnlyPaths=/etc/vefr /usr/local/bin/vefr

[Install]
WantedBy=multi-user.target
